# systemd/ai_wetter.service.example
[Unit]
Description=AI Wetter - Bulgarian Weather Forecast Tool
After=network.target

[Service]
User=your_username
Group=your_username
WorkingDirectory=/home/your_username/ai_wetter
Environment="PATH=/home/your_username/ai_wetter/venv/bin"
ExecStart=/home/your_username/ai_wetter/venv/bin/gunicorn 'app.app:create_app()' --bind=127.0.0.1:8000 --workers=2
Restart=always

[Install]
WantedBy=multi-user.target

# Replace 'your_username' with your actual username (e.g., 'jerry')